# extra/rc-scripts/gentoo/power-profile
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /home/cvsroot/gentoo-x86/sys-apps/speedfreq/files/speedfreq.rc,v 1.1 2003/08/21 02:31:55 vapier Exp $

depend() {
	need modules logger
	use speedfreq
}

start() {
	profile=$(/usr/local/sbin/hpdet power)
	ebegin "Selecting power profile ${profile}"
	/usr/local/sbin/hprofile power
	eend $?
}

